# atlas-service.conf
[program:atlas]
command = APP_PATH/gunicorn_start.sh
autostart=true
autorestart=true
stdout_logfile = APP_PATH/log/errors_atlas.log
redirect_stderr = true
